# root/usr/local/lib/systemd/system/lima-init.service.d/requires-cidata.conf
# kiwi-profile: lima

[Unit]
RequiresMountsFor=/mnt/lima-cidata
