# root/usr/local/lib/systemd/system/getty@.service.d/disable-on-wsl.conf
# Disable getty@ on WSL because it fails there.
[Unit]
ConditionVirtualization=!wsl
